# miros_ssh_proxy_1.0-1/DEBIAN/postinst
sudo systemctl enable miros-proxy-startup.service
chmod +x /usr/bin/miros/miros_ssh_proxy
