# miros_1.0-1/DEBIAN/postinst
sudo systemctl enable docker-startup.service
ln -s  /usr/bin/miros/miroslinuxamd64 /usr/bin/miros/miros
ln -s /usr/bin/miros/miros.desktop /usr/share/applications/miros.desktop
