# miros_1.0-1/DEBIAN/postinst
sudo systemctl enable docker-startup.service
