# miros_ssh_proxy_1.0-1/etc/systemd/system/miros-proxy-startup.service
[Unit]
Description=Runs Miros ssh proxy on startup

[Service]
Type=idle
ExecStart=/usr/bin/miros/miros_ssh_proxy

[Install]
WantedBy=multi-user.target
